# .config/zsh/profile.d/binaries
# User binaries #

path=( "$HOME/.local/bin" $path )

# vim: ft=zsh