# .config/zsh/profile.d/composer
# Composer binaries #

path=( "$HOME/.composer/vendor/bin" $path )

# vim: ft=zsh
