# .config/zsh/profile.d/makepkg
# Package manager user configuration #

export PACMAN='aura'

# vim: ft=zsh
