# .config/zsh/profile.d/pacman
# Package manager user configuration #

PACMAN='aura'

# vim: ft=zsh