# .config/zsh/profile.d/python
# Python environment #

export WORKON_HOME="$XDG_DATA_HOME/virtualenvs"

# vim: ft=zsh
