# .config/zsh/profile.d/scripts
# User shell scripts #

path=( "$XDG_DATA_HOME/zsh/scripts" $path )

# vim: ft=zsh
