# .config/zsh/profile.d/scripts
# User shell scripts #

path=( "$HOME/.local/scripts" $path )

# vim: ft=zsh